# templates/mcmyadmin_init.erb
#!/bin/bash
# /etc/init.d/mcmyadmin

### BEGIN INIT INFO
# Provides:   minecraft
# Required-Start: $local_fs $remote_fs
# Required-Stop:  $local_fs $remote_fs
# Should-Start:   $network
# Should-Stop:    $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description:    MCMyadmin Minecraft server
# Description:    Starts the MCMyadmin minecraft server
### END INIT INFO

#Settings
mcmyadmin_user="<%= @user %>"
mcmyadmin_session="mcmyadmin"
mcmyadmin_chdir="<%= @install_dir %>"
screen_cmd="/usr/bin/env screen"


ME=`whoami`
as_user() {
  if [ "$ME" == "${mcmyadmin_user}" ] ; then
    bash -c "$1"
  else
    su - ${mcmyadmin_user} -c "$1"
  fi
}

mc_start() {
  if mc_running; then
    echo "mcmyadmin is already running?"
    exit 1
  else
    as_user "cd ${mcmyadmin_chdir} && screen -dmS ${mcmyadmin_session} <%= @install_cmd %> <%= @mcma_run_args %>"
  fi
}

mc_stop() {
  if mc_running; then
    as_user "screen -p 0 -S ${mcmyadmin_session} -X eval 'stuff \"/quit\"\015'"
  else
    echo "mcmyadmin was not running."
    exit 1
  fi
}

mc_status() {
  if mc_running; then
    echo "mcmyadmin is running."
    exit 0
  else
    echo "mcmyadmin is not running."
    exit 1
  fi
}

mc_running() {
  check=`su -m ${mcmyadmin_user} -c "${screen_cmd} -ls" | grep ${mcmyadmin_session}`
  if [ "$check" ]; then
    return 0
  else
    return 1
  fi
}

#Start-Stop here
case "$1" in
  start)
    mc_start
  ;;
  stop)
    mc_stop
  ;;
  restart)
    mc_stop
    mc_start
  ;;
  status)
    mc_status
  ;;

  *)
    echo "Usage: /etc/init.d/mcmyadmin {start|stop|status|restart}"
    exit 1
  ;;
esac

exit 0
